# backend/PocketTube-ngrok.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  PocketTube-ngrok.xml
  ────────────────────
  Imports the "PocketTube-ngrok" scheduled task into Windows Task Scheduler.

  HOW TO IMPORT:
    1. Open Task Scheduler (search "Task Scheduler" in Start)
    2. Action menu → Import Task...
    3. Select this file
    4. Review and click OK
    - OR -
    Run from an admin CMD/PowerShell:
      schtasks /create /xml "PocketTube-ngrok.xml" /tn "PocketTube-ngrok"

  WHAT IT DOES:
    Runs start-ngrok.bat 5 seconds after user logon. The script itself then
    polls http://localhost:8080/health for up to 30 seconds before launching
    ngrok, so the 5-second scheduler delay + the health-check loop together
    guarantee the backend is ready before the tunnel opens.

  IMPORTANT — update the path below before importing:
    Replace  C:\Users\vibud\Desktop\02-01-Github\02-PocketTube\PocketTube\backend
    with the actual path to your backend\ folder if it ever moves.
-->
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Starts the PocketTube ngrok tunnel after the backend is healthy.</Description>
    <Author>PocketTube</Author>
  </RegistrationInfo>
  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
      <Delay>PT5S</Delay>
      <!-- 5-second delay after logon lets the backend task start first -->
    </LogonTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>C:\Users\vibud\Desktop\02-01-Github\02-PocketTube\PocketTube\backend\start-ngrok.bat</Command>
      <WorkingDirectory>C:\Users\vibud\Desktop\02-01-Github\02-PocketTube\PocketTube\backend</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
